# BCacheShutdownService_V1/bcache-stop.service
[Unit]
Description=Ensure clean shutdown of bcache devices
Before=umount.target
After=shutdown.target
DefaultDependencies=no

[Service]
ExecStart=__bcache_stop_script__
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=120
TimeoutStopSec=120

[Install]
WantedBy=umount.target
