# SystemdForDesktopUser/desktop.target
[Unit]
Description=Special target for user desktop session. Will be started after systemd env for user is tuned by special startup script.

